# deploy/sbom-backend.service
[Unit]
Description=SBOM Platform Backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/sbom/backend
EnvironmentFile=/var/www/sbom/.env
Environment=STATIC_DIR=/var/www/sbom/frontend/dist

# 1 worker: 1GB RAM 伺服器的安全設定
ExecStart=/var/www/sbom/backend/venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 9100 \
    --workers 1 \
    --no-access-log \
    --timeout-keep-alive 30

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 記憶體上限保護：超過 400MB 自動重啟
MemoryMax=400M

[Install]
WantedBy=multi-user.target
